# dist/debian/mariner3d.marinervideo.service
[Unit]
Description=Mariner3d Video
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python /opt/venvs/mariner3d/lib/python3.9/site-packages/mariner/server/video.py
Restart=on-abort

[Install]
WantedBy=multi-user.target